# roles/prototypes/provision_apache/files/etc/cron.d/mrtg
*/5     *       *       *       *       root    env LANG=C /usr/local/bin/mrtg  /usr/local/www/apache24/htdocs/mrtg/mrtg.cfg
